# docker/app/supervisor.conf
[program:php-fpm]
command=php-fpm -R -F -c /etc/php/7.1/fpm/php-fpm.conf

[program:nginx]
command=nginx

[program:cron]
command=cron -f
