# docker/cron/crontab
11 * * * * . /root/project_env.sh; /usr/local/bin/php /var/www/html/artisan import:all > /proc/1/fd/1 2>/proc/1/fd/2
